# supervisor.cfg
[program:vadd]
command=/home/work/online/bin/python main.py -debug=False -port=96%(process_num)02d
process_name=%(program_name)s-99%(process_num)02d
numprocs=5
directory=/home/work/online/src/vadd_ws/
environment=PATH="/home/work/online/bin"
startretries=2
startsecs=5
autostart=yes
autorestart=true
redirect_stderr=true
logfile_maxbytes=100M
stdout_logfile=/data/super_log/vadd_ws.log

[group:vadd_ws]
programs=vadd_ws
priority=999